# setup_agv.service
[Unit]
Description=Setup AGV Ethernet Connection
After=network-online.target

[Service]
ExecStart=/YOUR_PATH/setup_agv.sh
RemainAfterExit=yes
Restart=always
User=root

[Install]
WantedBy=multi-user.target